# etc/xdg/autostart/x11vnc.desktop
#bayram KARAHAN bayramk@gmail.com
[Desktop Entry]
Exec=sh /usr/bin/x11vncusers.sh
